# templates/default/etc/init.d/browsermob.erb
#!/bin/sh

#
# chkconfig: - 80 20
# processname: browsermob
# description: Browsermob Proxy
#

### BEGIN INIT INFO
# Provides: browsermob
# Required-Start: $network $syslog
# Required-Stop: $network $syslog
# Default-Start:
# Default-Stop:
# Short-Description: start and stop browsermob
### END INIT INFO

. /etc/rc.d/init.d/functions

DAEMON="<%= node['aet']['browsermob']['root_dir']  %>/current/bin/browsermob-proxy"
PARAM="-port <%= node['aet']['browsermob']['port'] %> -proxyPortRange <%= node['aet']['browsermob']['proxy_port_range'] %>"
LOCK_FILE="/var/lock/subsys/browsermob"
USER="<%= node['aet']['browsermob']['user'] %>"

browsermob_pid()
{
  echo `ps -aefw | grep 'browsermob/current' | grep -v grep | awk '{print $2}'`
}

do_start()
{
  if [ ! -f "$LOCK_FILE" ] ; then
    echo -n "Starting browsermob: "
    runuser -l "$USER" -c "$DAEMON $PARAM >> <%= node['aet']['browsermob']['log_dir'] %>/browsermob.log 2>&1 &" && echo_success || echo_failure
    RETVAL=$?
    echo
    if [ $RETVAL -eq 0 ] ; then
      touch $LOCK_FILE
    else
      echo "Error while starting $DAEMON. Error code from browsermob: [$RETVAL]"
    fi
  else
    echo "$DAEMON is locked. Lock file [$LOCK_FILE] is present therefore instance is probably running."
    RETVAL=0
  fi
}
do_stop()
{
  echo -n "Stopping browsermob: "
  pid=$(browsermob_pid)
  kill $pid > /dev/null 2>&1 && echo_success || echo_failure
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && rm -f $LOCK_FILE
}
do_status()
{
  pid=$(browsermob_pid)
  if [ -n "$pid" ]
  then
     echo "Browsermob is running with pid: $pid"
     RETVAL=0
  else
     echo "Browsermob is not running"
     RETVAL=3
  fi
}

case "$1" in
  start)
    do_start
  ;;
  stop)
    do_stop
  ;;
  restart)
    do_stop
    do_start
  ;;
  status)
    do_status
  ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    RETVAL=1
esac

exit $RETVAL
